# asterisk4ucs-locallistener/debian/postinst
#!/bin/sh

invoke-rc.d univention-directory-listener restart

